# bin/crypt.service
[Unit]
Description=Encryption

[Service]
ExecStart= /home/client/src/cryptography/bin/crypt.sh

[Install]
WantedBy=multi-user.target
