# bin/crypt.service
[Unit]
Description=Decryption

[Service]
ExecStart= /home/client/bin/crypt.sh

[Install]
WantedBy=multi-user.target
